# stock.cron
* * * * * cd /Users/hello/wwwroot/stock && php artisan schedule:run >>  /Users/hello/wwwroot/stock/log 2>&1
