# stock.cron
* * * * * cd /wwwroot/stock && php artisan schedule:run >>  /wwwroot/storage/logs 2>&1
0 1 * * 6 cd /wwwroot/stock && php artisan StockFlowSpider yes
0 16 * * 6 cd /wwwroot/stock && php artisan StockAnalyzer yes
0 23 * * 6 cd /wwwroot/stock && php artisan MACD yes
0 4 * * 0 cd /wwwroot/stock && php artisan KDJ yes
