# stock.cron
* * * * * cd /wwwroot/stock/ && php artisan schedule:run >>  /wwwroot/log 2>&1
0 8 * * 7 cd /wwwroot/stock/ && php artisan MACD yes
0 1 * * 7 cd /wwwroot/stock/ && php artisan KDJ yes
0 21 * * 5 cd /wwwroot/stock/ && php artisan StockFlowSpider yes
0 14 * * 6 cd /wwwroot/stock/ && php artisan StockAnalyzer yes


# 添加crontab :  crontab stock.cron